# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:glusterd]
command=/usr/sbin/glusterd -p /var/run/glusterd.pid 