# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -f /etc/rsyslog.conf -d

[program:glusterd]
command=/usr/sbin/glusterd -p /var/run/glusterd.pid 
